# scripts/systemd/rebar-mgmt.service
[Unit]
Description=钢筋精细化管理平台 Docker Compose 启动服务
Documentation=https://github.com/rebar-mgmt/rebar-mgmt
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/rebar-mgmt
ExecStart=/opt/rebar-mgmt/scripts/rebar-startup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
